# ops/launchd/com.discordclone.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.discordclone.frontend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/mohamed.darwish/Documents/New project' && npm run dev -- --host 127.0.0.1 --port 5173</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mohamed.darwish/Documents/New project</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/discordclone-frontend.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/discordclone-frontend.err.log</string>
</dict>
</plist>
